# rp-auto-program@.service
[Unit]
Description=Program RP device attached to '%I'

[Service]
User=<user>
ExecStart=<path>/rp-program-shim.sh '%I'
